# res/systemd/_bthelper_overide.service
[Service]
Type=oneshot
ExecStartPost=/bin/hciconfig %I piscan
ExecStartPost=/bin/hciconfig %I sspmode 1